# templates/netronome/netronom-agent.service.j2
[Unit]
Description=Netronome vnstat Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=netronome
Group=netronome
ExecStart=/usr/local/bin/netronome agent --config /home/{{ ansible_user }}/.config/netronome/agent.toml
Restart=always
RestartSec=10

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/netronome

[Install]
WantedBy=multi-user.target
